# claude-window-primer.service
[Unit]
Description=Claude Pro 5-hour window primer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/testing/claude-window-primer
ExecStart=%h/testing/claude-window-primer/primer.py bot
Environment=PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
